# src/ubuntu/base/20.04-xfce/supervisor/vncserver.conf


[program:vncserver]
command=/opt/supervisor-scripts/launch_vnc.sh
user=desktopus
autostart=true
stderr_logfile=/var/log/vncserver.err.log
stdout_logfile=/var/log/vncserver.out.log
priority=1000